# templates/raft-observer.timer.j2
[Unit]
Description=Patroni Raft observer timer

[Timer]
OnActiveSec=10s
OnUnitInactiveSec=10s
AccuracySec=1s

[Install]
WantedBy=timers.target
